# python-3.5/supervisor/django.ini
[program:toligado]
process_name = %(program_name)s_%(process_num)02d
command = /usr/local/bin/uwsgi --chdir /var/www/html/dev/ -s /tmp/uwsgi.sock --master --workers 1 --logto /var/www/html/dev/supervisord/logs/uwsgi.log --module toligado.wsgi:application --manage-script-name --mount /=toligado.wsgi:application --env DJANGO_SETTINGS_MODULE=toligado.settings
directory = /var/www/html/dev/
pythonpath = /var/www/html/dev/
chmod-socket = 666
socket = /tmp/uwsgi.sock
user = root
uid = root
gid = root
master = True
harakiri = 20
numprocs = 4
autostart = true
autorestart = true
redirect_stderr = true
stopwaitsecs = 60
stopsignal = INT